# systemd/vnstat-web-update.service
[Unit]
Description=Update vnstat web panel data

[Service]
Type=oneshot
ExecStart=/usr/local/bin/vnstat-web-update.sh
