# systemd/vnstat-web-update.service
[Unit]
Description=Generate vnStat web data (json/png/summary)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/vnstat-web-update.sh
